# rootfs/var/spool/cron/crontabs/www-data
*/5 * * * * php -f /var/www/nextcloud/cron.php
